# deploy/jarvis.service
[Unit]
Description=JARVIS MK1 Lite - Telegram Interface for Claude Code
After=network.target

[Service]
Type=simple
# Run as root for maximum Claude Code autonomy
User=root
Group=root
WorkingDirectory=/opt/jarvis-mk1
Environment=PATH=/opt/jarvis-mk1/venv/bin:/usr/local/bin:/usr/bin:/root/.local/bin
ExecStart=/opt/jarvis-mk1/venv/bin/python -m jarvis_mk1_lite
Restart=always
RestartSec=10

# No restrictions - full system access for Claude Code
# NoNewPrivileges=false (default)
# ProtectSystem=false (default)
# ProtectHome=false (default)

[Install]
WantedBy=multi-user.target
